# util/koji-sweep-db.timer
[Unit]
Description=Daily maintenance script for koji db
Documentation=https://pagure.io/docs/koji/

[Timer]
OnCalendar=daily
Accuracy=1h
Persistent=true

[Install]
WantedBy=timers.target
